# etc/cron.hourly/sigi
#!/bin/bash
cd /srv/sigi/sigi
python manage.py gera_map_data
