# etc/cron.hourly/sigi
#!/bin/bash
cd /var/aplicacoes/django/sigi/sigi
python manage.py gera_map_data
