# etc/supervisor/conf.d/sigi.conf
[program:sigi]
command = /srv/sigi/bin/run_sigi
user = sigi
stdout_logfile = /var/log/sigi/sigi-supervisor.log
redirect_stderr = true
