# etc/cron.daily/sigi
#!/bin/bash
cd /srv/sigi/sigi/
python manage.py atualiza_uso_servico -v 0
python manage.py clearsessions -v 0
